# runit/notifier-service/run
#!/bin/sh

exec 2>&1
DISPLAY=':0' exec /usr/bin/notifier
